# airootfs/etc/udev/rules.d/91-apple-superdrive.rules
ACTION=="add", SUBSYSTEM=="block", ATTRS{idVendor}=="05ac", ATTRS{idProduct}=="1500", RUN+="/usr/bin/sg_raw --cmdset=1 /dev/%k EA 00 00 00 00 00 01"
